# automation/insights/com.dneves.insights-nextmeal.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  The DAILY next-meal job — 17:00 local, before dinner. Assembles 3 plates from the
  day's remaining budget and writes them to the `next_meal` cache tab.

  Installed copy at ~/Library/LaunchAgents/com.dneves.insights-nextmeal.plist:

    cp automation/insights/com.dneves.insights-nextmeal.plist ~/Library/LaunchAgents/
    launchctl bootout  gui/$(id -u)/com.dneves.insights-nextmeal 2>/dev/null
    launchctl bootstrap gui/$(id -u) ~/Library/LaunchAgents/com.dneves.insights-nextmeal.plist

  17:00 is chosen so the suggestion reflects what was actually eaten through the day
  and lands while dinner is still a decision (the card is stamped with this time, since
  a late-logged lunch would make it slightly stale — an accepted trade of local
  generation). Fill INGEST_TOKEN / HEALTH_BACKEND_URL as in the weekly plist.
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.dneves.insights-nextmeal</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/dneves/personal/health-tracker/backend/venv/bin/python</string>
        <string>/Users/dneves/personal/health-tracker/automation/insights/generate.py</string>
        <string>next-meal</string>
    </array>

    <!-- Every day at 17:00 local time. -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key><integer>17</integer>
        <key>Minute</key><integer>0</integer>
    </dict>

    <key>RunAtLoad</key>
    <false/>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Users/dneves/.nvm/versions/node/v22.12.0/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>HEALTH_BACKEND_URL</key>
        <string>REPLACE_WITH_CLOUD_RUN_URL</string>
        <key>INGEST_TOKEN</key>
        <string>REPLACE_WITH_APP_TOKEN</string>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/dneves/personal/health-tracker/automation/insights/logs/launchd.out.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/dneves/personal/health-tracker/automation/insights/logs/launchd.err.log</string>
</dict>
</plist>
