# automation/insights/com.dneves.insights-weekly.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  The WEEKLY insights job — Sunday 09:00 local. Reviews the 7 days that just ended and
  writes the coaching report to the `weekly_reports` tab.

  This file lives in the repo; a COPY is installed at
  ~/Library/LaunchAgents/com.dneves.insights-weekly.plist (launchd only reads there):

    cp automation/insights/com.dneves.insights-weekly.plist ~/Library/LaunchAgents/
    launchctl bootout  gui/$(id -u)/com.dneves.insights-weekly 2>/dev/null
    launchctl bootstrap gui/$(id -u) ~/Library/LaunchAgents/com.dneves.insights-weekly.plist

  Fill INGEST_TOKEN and HEALTH_BACKEND_URL below (the app's X-Auth-Token and the Cloud
  Run base URL). RunAtLoad is OFF so bootstrapping never fires an unexpected write; a
  missed Sunday (Mac asleep) is caught up shortly after wake.
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.dneves.insights-weekly</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/dneves/personal/health-tracker/backend/venv/bin/python</string>
        <string>/Users/dneves/personal/health-tracker/automation/insights/generate.py</string>
        <string>weekly</string>
    </array>

    <!-- Sunday (Weekday 0) at 09:00 local time. -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Weekday</key><integer>0</integer>
        <key>Hour</key><integer>9</integer>
        <key>Minute</key><integer>0</integer>
    </dict>

    <key>RunAtLoad</key>
    <false/>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Users/dneves/.nvm/versions/node/v22.12.0/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>HEALTH_BACKEND_URL</key>
        <string>REPLACE_WITH_CLOUD_RUN_URL</string>
        <key>INGEST_TOKEN</key>
        <string>REPLACE_WITH_APP_TOKEN</string>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/dneves/personal/health-tracker/automation/insights/logs/launchd.out.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/dneves/personal/health-tracker/automation/insights/logs/launchd.err.log</string>
</dict>
</plist>
